# supervisord.conf

[supervisord]
logfile = /Users/biajee/Code/supervisor/logs/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = /tmp/supervisord.pid
nodaemon = False
minfds = 1024
minprocs = 200
umask = 022
identifier = supervisor
directory = /tmp
nocleanup = true
childlogdir = /tmp

[program:flaskbyexample]
command = gunicorn app:app -b localhost:8000
directory = "/Users/biajee/Code/flask-by-example"
user = biajee
